# servers/gconode01/files/etc/supervisor/conf.d/adcp.conf
[program:collect_adcp]
command=/home/gcouser/dev/uwobs/common/apps/collectors/collect_adcp.sh
user=gcouser
environment=HOME="/home/gcouser",USER="gcouser"
autorestart=true
stopasgroup=true
