# servers/gconode01/files/etc/supervisor/conf.d/ctd.conf
[program:collect_ctd]
command=/home/gcouser/bin/collect_ctd.sh
user=gcouser
environment=HOME="/home/gcouser",USER="gcouser"
stopasgroup=true
