# servers/gconode01/files/etc/supervisor/conf.d/fluorometer.conf
[program:collect_fluorometer]
command=/home/gcouser/bin/collect_fluorometer.sh
user=gcouser
environment=HOME="/home/gcouser",USER="gcouser"
stopasgroup=true
