# servers/gconode01/files/etc/supervisor/conf.d/vemco.conf
[program:collect_vemco]
command=/home/gcouser/bin/collect_vemco.sh
user=gcouser
environment=HOME="/home/gcouser",USER="gcouser"
stopasgroup=true
